# deploy/linux/brick-scanner-catalog-refresh.timer
# Fires brick-scanner-catalog-refresh.service daily at 07:30 local time,
# shortly after Rebrickable's ~07:12 ET catalog update (matches the launchd
# StartCalendarInterval Hour=7 Minute=30).
#
# Persistent=true: if the VM was powered off at 07:30, the job runs at the next
# boot — the systemd analogue of launchd "runs at next wake".

[Unit]
Description=Daily Brick Scanner catalog refresh

[Timer]
OnCalendar=*-*-* 07:30:00
Persistent=true

[Install]
WantedBy=timers.target
