# deploy/linux/brick-scanner-minifig-prices.service
# Brick Scanner — daily BrickLink price refresh for the My Minifigs collection
# (systemd oneshot).
#
# Linux equivalent of com.brickscanner.minifig-prices.plist. Run on a schedule
# by brick-scanner-minifig-prices.timer (05:00 daily). LOCAL-ONLY (does nothing
# useful on Render). The script writes its own minifig_prices.log; stray
# stdout/stderr goes to the journal:
#   journalctl -u brick-scanner-minifig-prices
#
# TEMPLATE: __PROJECT_DIR__, __USER__, __PY__ filled in by ./install_systemd.sh.

[Unit]
Description=Brick Scanner minifig BrickLink price refresh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=__USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PY__ refresh_minifig_prices.py
