# deploy/linux/brick-scanner-minifig-prices.timer
# Fires brick-scanner-minifig-prices.service daily at 05:00 local time
# (matches the launchd StartCalendarInterval Hour=5 Minute=0).
#
# Persistent=true: if the VM was off at 05:00, the job runs at the next boot.

[Unit]
Description=Daily Brick Scanner minifig price refresh

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true

[Install]
WantedBy=timers.target
